# infra/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Cardsapp
After=docker.service
Requires=docker.service

[Service]
Restart=always
RestartSec=5
ExecStartPre=-/usr/bin/docker rm -f cloudflare-tunnel
ExecStart=/bin/sh -c '. ~/server-project/.env && exec docker run --rm \
  --name cloudflare-tunnel \
  --network cardsapp-net \
  -e TUNNEL_TOKEN="$CLOUDFLARE_TOKEN" \
  cloudflare/cloudflared:latest \
  tunnel --no-autoupdate run'
ExecStop=/usr/bin/docker stop cloudflare-tunnel

[Install]
WantedBy=multi-user.target
